# node-react-jssor v7/kiosk-ipray-setup/new-files/home-ipray-config/autostart/gnome-keyring-secrets.desktop
[Desktop Entry]
Type=Application
Name=Secret Storage Service
Comment=GNOME Keyring: Secret Service (disabled for kiosk autologin)
Exec=/bin/true
Hidden=true
X-GNOME-Autostart-enabled=false
